# net-misc/sunloginclient/files/runoraydaemon.initd
#!/sbin/openrc-run
# Copyright 1999-2021 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2

description="check running sunloginclinet loop"
command="/opt/sunlogin/bin/oray_rundaemon"
command_args="-m server"
command_background=true
pidfile="/run/${RC_SVCNAME}.pid"

depend() {
	need net
	before runsunloginclient
}
